# files/service/graphite-stresser/run
#!/bin/sh
set +x
# Setup defaults if needed
if [ -z ${GRAPHITE_STRESSER_HOST} ]; then
    GRAPHITE_STRESSER_HOST=localhost   
fi

if [ -z ${GRAPHITE_STRESSER_PORT} ]; then
    GRAPHITE_STRESSER_PORT=2003 
fi
if [ -z ${GRAPHITE_STRESSER_NUMHOSTS} ]; then
    GRAPHITE_STRESSER_NUMHOSTS=1   
fi
if [ -z ${GRAPHITE_STRESSER_NUMTIMERS} ]; then
    GRAPHITE_STRESSER_NUMTIMERS=128   
fi
if [ -z ${GRAPHITE_STRESSER_INTERVAL} ]; then
    GRAPHITE_STRESSER_INTERVAL=10  
fi
if [ -z ${GRAPHITE_STRESSER_DEBUG} ]; then
    GRAPHITE_STRESSER_DEBUG=false 
fi
if [ ! -z ${GRAPHITE_STRESSER_LOGFILE} ]; then
    GRAPHITE_STRESSER_LOGGING=" >> /var/log/graphite-stresser.log 2>&1"
fi

echo "Running graphite-stresser ${GRAPHITE_STRESSER_PORT} ${GRAPHITE_STRESSER_NUMHOSTS} ${GRAPHITE_STRESSER_NUMTIMERS} ${GRAPHITE_STRESSER_INTERVAL} ${GRAPHITE_STRESSER_DEBUG}"

exec java -jar ${WORKDIR}/dist/stresser.jar \
    ${GRAPHITE_STRESSER_HOST} \
    ${GRAPHITE_STRESSER_PORT} \
    ${GRAPHITE_STRESSER_NUMHOSTS} \
    ${GRAPHITE_STRESSER_NUMTIMERS} \
    ${GRAPHITE_STRESSER_INTERVAL} \
    ${GRAPHITE_STRESSER_DEBUG} \
    ${GRAPHITE_STRESSER_LOGGING}

