# files/service/graphite-stresser/run
#!/bin/bash
set -x

# DEFAULT SCHEDULE SETTINGS 
if [ -z ${GRAPHITE_STRESSER_HOST} ]; then
    GRAPHITE_STRESSER_HOST=localhost   
fi

if [ -z ${GRAPHITE_STRESSER_PORT} ]; then
    GRAPHITE_STRESSER_PORT=2003 
fi
if [ -z ${GRAPHITE_STRESSER_NUMHOSTS} ]; then
    GRAPHITE_STRESSER_NUMHOSTS=1   
fi
if [ -z ${GRAPHITE_STRESSER_NUMTIMERS} ]; then
    GRAPHITE_STRESSER_NUMTIMERS=128   
fi
if [ -z ${GRAPHITE_STRESSER_INTERVAL} ]; then
    GRAPHITE_STRESSER_INTERVAL=10  
fi
if [ -z ${GRAPHITE_STRESSER_DEBUG} ]; then
    GRAPHITE_STRESSER_DEBUG=false 
fi
if [ ! -z ${GRAPHITE_STRESSER_LOGFILE} ]; then
    GRAPHITE_STRESSER_LOGGING=" >> ${GRAPHITE_STRESSER_LOGFILE} 2>&1"
fi

# DEFAULT ACCUMULATION SETTINGS
if [ -z ${GRAPHITE_STRESSER_METRICS} ]; then
    GRAPHITE_STRESSER_METRICS=10000
fi
if [ -z ${GRAPHITE_STRESSER_ACCUMULATION_INTERVAL} ]; then
    GRAPHITE_STRESSER_ACCUMULATION_INTERVAL=10
fi
if [ -z ${GRAPHITE_STRESSER_JITTER} ]; then
    GRAPHITE_STRESSER_JITTER=10
fi
if [ -z ${GRAPHITE_STRESSER_AGENTS} ]; then
    GRAPHITE_STRESSER_AGENTS=100
fi

if [ "x${GRAPHITE_STRESSER_MODE}" == "xscheduled" ]; then
	echo "Running graphite-stresser ${GRAPHITE_STRESSER_PORT} ${GRAPHITE_STRESSER_NUMHOSTS} ${GRAPHITE_STRESSER_NUMTIMERS} ${GRAPHITE_STRESSER_INTERVAL} ${GRAPHITE_STRESSER_DEBUG}"
	
	exec java -jar ${WORKDIR}/dist/stresser.jar \
	    ${GRAPHITE_STRESSER_HOST} \
	    ${GRAPHITE_STRESSER_PORT} \
	    ${GRAPHITE_STRESSER_NUMHOSTS} \
	    ${GRAPHITE_STRESSER_NUMTIMERS} \
	    ${GRAPHITE_STRESSER_INTERVAL} \
	    ${GRAPHITE_STRESSER_DEBUG} \
	    ${GRAPHITE_STRESSER_LOGGING}
	    
elif [ "x${GRAPHITE_STRESSER_MODE}" == "xaccumulate" ]; then
	echo "Running haggar ${GRAPHITE_STRESSER_PORT} ${GRAPHITE_STRESSER_NUMHOSTS} ${GRAPHITE_STRESSER_NUMTIMERS} ${GRAPHITE_STRESSER_INTERVAL} ${GRAPHITE_STRESSER_DEBUG}"
	
	exec ${WORKDIR}/go/bin/haggar \
  		--prefix="STRESS_ACCUMULATION" \
	    --agents=${GRAPHITE_STRESSER_AGENTS} \
	    --carbon=${GRAPHITE_STRESSER_HOST}:${GRAPHITE_STRESSER_PORT} \
	    --metrics=${GRAPHITE_STRESSER_METRICS} \
	    --flush-interval=${GRAPHITE_STRESSER_INTERVAL}s \
	    --spawn-interval=${GRAPHITE_STRESSER_ACCUMULATION_INTERVAL}s \
	    --jitter=${GRAPHITE_STRESSER_JITTER}s
else
	echo graphite-stresser: Set GRAPHITE_STRESSER_MODE equal to 'scheduled' or 'accumulate'
	exit 1
fi
